# root/etc/s6-overlay/s6-rc.d/svc-minio/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc \
	/app/minio \
	server \
	/config \
	--console-address ":9001"
